# dialer-backend.service
[Unit]
Description=AK Dialer Backend API
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/AK_Dialer/backend
Environment="PATH=/root/AK_Dialer/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
# Database Configuration
Environment="DATABASE_URL=postgresql://your_username:your_password@localhost:5432/your_database"
# Security
Environment="SECRET_KEY=your-secret-key-change-in-production"
Environment="ALGORITHM=HS256"
Environment="ACCESS_TOKEN_EXPIRE_MINUTES=1440"
# Asterisk Configuration
Environment="USE_MOCK_DIALER=False"
Environment="ASTERISK_HOST=163.245.208.168"
Environment="ASTERISK_AMI_PORT=5038"
Environment="ASTERISK_AMI_USERNAME=admin"
Environment="ASTERISK_AMI_PASSWORD=amp111"
Environment="ASTERISK_CONTEXT=from-internal"
Environment="ASTERISK_TRUNK=PJSIP/trunk"
# CORS
Environment="CORS_ORIGINS=http://163.245.208.168:3000,http://localhost:3000"
# Call Settings
Environment="DEFAULT_CALL_TIMEOUT=60"
Environment="MAX_CONCURRENT_CALLS=10"
# Wait for PostgreSQL to be ready before starting
ExecStartPre=/bin/bash -c 'until pg_isready -h localhost -p 5432; do sleep 1; done'
# Increase startup timeout to allow AMI connection
TimeoutStartSec=300
TimeoutStopSec=30
ExecStart=/root/AK_Dialer/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
